# .config/systemd/user/rclone-gdrive.timer
[Unit]
Description=RClone GDrive Sync timer

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h
Unit=rclone-gdrive.service

[Install]
WantedBy=timers.target
