# .config/systemd/user/rclone-gdrive.service
[Unit]
Description=RClone GDrive Sync

[Service]
Type=simple
ExecStart=~/.bin/rclone-gdrive
