# .config/systemd/user/sync-gdrive.timer
[Unit]
Description=Sync Google Drive timer

[Timer]
OnBootSec=1d
OnUnitActiveSec=1d
Unit=sync-gdrive.service

[Install]
WantedBy=timers.target
