# .config/systemd/user/sync-gdrive.service
[Unit]
Description=Sync Google Drive

[Service]
Type=simple
ExecStart=~/.bin/sync-gdrive
